# linux/tmux/crontab
@reboot . /home/sen4ik/.profile ; sh /home/sen4ik/startTmux.sh &>> log
